# roles/observium/files/letsencrypt.timer
[Unit]
Description=Daily Letsencrypt certificate renewal

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
