# roles/observium/files/mysqlcheck.service
[Unit]
Description=Optimize all MySQL databases
Requires=mariadb.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mysqlcheck -Aos
